# mongo/service/mongod-shard/run
#!/bin/bash
	ulimit -SHn 10240
	exec 2>&1
	exec /usr/bin/mongod -f /etc/mongod.conf --shardsvr --dbpath /data/mongodb/data
